# crontab
10 */4 * * * bash /scripts/gcs-storage-sync.bash
